# minecraft/stop-minecraft.service.j2
[Unit]
Description=Service for stopping minecraft server
[Service]
Type=oneshot
ExecStart=/bin/systemctl stop minecraft.service
