# minecraft/stop-papermc.service.j2
[Unit]
Description=Service for stopping papermc server
[Service]
Type=oneshot
ExecStart=/bin/systemctl stop papermc.service
